# Sources/EkkoPlatform/Resources/com.ekko.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.ekko.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>__CLI_PATH__</string>
        <string>--agent-trigger</string>
    </array>

    <key>StartCalendarInterval</key>
    __SCHEDULE__

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__LOG_DIR__/agent.log</string>

    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/agent-error.log</string>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
